# private_dot_config/systemd/user/argocd-monitor-oneshot.service
[Unit]
Description=ArgoCD Event Monitor Agent (One-shot)
Documentation=https://github.com/laurigates/dotfiles
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/argocd-monitor --once

# Environment file for secrets
EnvironmentFile=-%h/.config/argocd-monitor/env

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
